# systemd/driveforge-daemon.service
[Unit]
Description=DriveForge daemon (orchestration, web UI, REST API)
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=driveforge
Group=driveforge
ExecStart=/usr/bin/driveforge-daemon
Restart=on-failure
RestartSec=5
WorkingDirectory=/var/lib/driveforge
StateDirectory=driveforge
LogsDirectory=driveforge

# Drive access — the daemon needs privileged access to block devices for
# SMART reads, secure-erase, and badblocks. Keep the blast radius contained
# as much as systemd allows.
#
# CAP_SETUID + CAP_SETGID are in the bounding set (NOT ambient) so that
# the v0.3.1 in-app self-update can sudo into systemctl. Sudo is setuid-
# root and needs to call setresuid(0) / setresgid(0) — those calls require
# CAP_SETUID / CAP_SETGID to be present in the bounding set, otherwise
# they're dropped permanently for the entire process tree even through
# setuid execution. Bug shipped in v0.3.1 (in-app update never exercised
# on the real protected daemon — the test suite monkeypatches subprocess);
# fixed in v0.4.2.
#
# Trade-off: an attacker who compromises the daemon could exec a setuid-
# root binary and elevate. That capability is the same one sudo needs to
# function at all — the alternative would be to drop sudo entirely and
# use a polkit-mediated D-Bus call to systemd's StartUnit method, which
# is a meaningful refactor. For a homelab tool with no network auth in
# the first place, the trade-off is acceptable.
AmbientCapabilities=CAP_SYS_RAWIO CAP_SYS_ADMIN CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_SYS_ADMIN CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_AUDIT_WRITE

NoNewPrivileges=false
ProtectHome=true
ProtectSystem=strict
PrivateTmp=true
# /etc/hostname + /etc/hosts are needed for the v0.2.8+ Settings → Hostname
# rename feature. Without these, the daemon hits EROFS (Errno 30) trying
# to write to /etc — systemd's mount-namespace protection makes the whole
# filesystem read-only to the daemon by default, even though shells see
# normal write permissions. Bug shipped in v0.2.8 (rename feature never
# exercised against the real protected daemon — dev mode short-circuits
# the writes); fixed in v0.4.1.
ReadWritePaths=/var/lib/driveforge /var/log/driveforge /etc/driveforge /etc/hostname /etc/hosts

# Allow journald + local logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
